# wifi-rssi/kg.service
[Unit]
Description=LED Flash

[Service]
Type=simple
ExecStart=/bin/sh -c /home/root/wifiRssi.sh
